# config/includes.chroot/etc/systemd/system/sanctum-flatpak-update.timer
[Unit]
Description=Daily automatic Flatpak updates (Firefox, Telegram)

[Timer]
OnBootSec=10min
OnCalendar=daily
RandomizedDelaySec=1h
# Persistent only affects OnCalendar= triggers — hence the calendar form.
Persistent=true

[Install]
WantedBy=timers.target
